# infra/catastro-api.service
[Unit]
Description=Catastro Chile API (FastAPI/Uvicorn)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/catastral.cl/backend
EnvironmentFile=/var/www/catastral.cl/.env
ExecStart=/var/www/catastral.cl/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
